# selmabot/selma.service
[Unit]
Description=Selma Bot
After=syslog.target

[Service]
Restart=always
Type=simple
WorkingDirectory=/home/bot/bot
ExecStart=python3 /home/bot/bot/main.py
SyslogIdentifier=selma
RestartSec=5

[Install]
WantedBy=multi-user.target
